# scripts/quant-data-collector.service
[Unit]
Description=Quant Research Lab - Data Collector Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=bot2
Group=bot2
WorkingDirectory=/home/bot2/claude/quant_research_lab
Environment="PATH=/usr/bin:/usr/local/bin"
EnvironmentFile=/home/bot2/claude/quant_research_lab/.env

# Run data collection every hour
ExecStart=/usr/bin/python3 -c "
import time
import schedule
from datetime import datetime, timedelta
from data.market_data_collector import MarketDataCollector
from data.mysql_storage import MySQLStorage
import os

def collect_data():
    print(f'[{datetime.now()}] Starting data collection...')
    storage = MySQLStorage(
        host=os.getenv('MYSQL_HOST', 'localhost'),
        user=os.getenv('MYSQL_USER', 'bot2'),
        password=os.getenv('MYSQL_PASSWORD', 'bot2_123456'),
        database=os.getenv('MYSQL_DATABASE', 'quant_research')
    )
    storage.connect()
    collector = MarketDataCollector(storage=storage, symbols=['ETHUSDT'])
    end = datetime.utcnow()
    start = end - timedelta(hours=2)
    collector.download_all_historical(start_date=start, end_date=end)
    collector.close()
    storage.disconnect()
    print(f'[{datetime.now()}] Data collection complete.')

schedule.every().hour.do(collect_data)

while True:
    schedule.run_pending()
    time.sleep(60)
"

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=quant-data-collector

[Install]
WantedBy=multi-user.target
